# deb/asl-update-node-list_2.0.0-beta.4_all/DEBIAN/prerm
#!/bin/sh -e

systemctl disable update-node-list
systemctl stop update-node-list
